# packages/usr/autostart/com.github.wwmm.easyeffects.desktop
[Desktop Entry]
Comment=Easy Effects Service
Exec=easyeffects --hide-window --service-mode
Icon=com.github.wwmm.easyeffects
Name=Easy Effects
StartupNotify=false
Terminal=false
Type=Application
X-GNOME-Autostart-Phase=Application
X-KDE-autostart-phase=2
